# Dockerfile/gobblin/supervisor/gobblin.conf
[program:gobblin]
startretries=0
command=/opt/gobblin/gobblin-dist/bin/gobblin-standalone.sh start
autostart=true
autorestart=false
startsecs=0